# wl11000/Image_static/etc/inittab
ttyS0::sysinit:/etc/rc S
ttyS0::wait:/etc/rc
ttyS0::askfirst:/bin/ash -s
