# Image_static/etc/inittab
console::wait:/etc/rc S
console::wait:/etc/rc
console::askfirst:/bin/ash -s
